# workstation-container/c9.conf
[program:c9]
command="node server.js -a : -w /dev --port 9999"
cwd=/usr/local/c9/core
autostart=true
autorestart=true