# workstation-container/jupyter.conf
[program:notebook]
command=/usr/local/bin/tini -- /usr/local/bin/start-notebook.sh
directory=/root
autostart=true
autorestart=true