# base-container/jupyter.conf
[program:notebook]
command="/usr/local/bin/start-notebook.sh"
cwd=/usr/local/c9/core
autostart=true
autorestart=true